# profile.d/91_pager
export PAGER=/usr/bin/less
export LESS="${LESS//[r]/}r"
