# less/.config/profile.d/configure-less-pager.sh
#!/bin/sh
export PAGER=/usr/bin/less
LESS=$(echo "${LESS}" | tr -d RM)
export LESS="${LESS}RM"
