# less/.config/profile.d/configure-less-pager.sh
export PAGER=/usr/bin/less
export LESS="${LESS//[RM]/}RM"
