# base/.config/profile.d/source-xdg-user-dirs.sh
XDG_CONFIG_HOME=$HOME/etc
source ${XDG_CONFIG_HOME}/user-home.dirs
source ${XDG_CONFIG_HOME}/user-dirs.dirs
