# base/.config/profile.d/source-xdg-user-dirs.sh
#!/bin/sh
XDG_CONFIG_HOME=$HOME/.config
source ${XDG_CONFIG_HOME}/user-home.dirs
source ${XDG_CONFIG_HOME}/user-dirs.dirs
